# real lab/start_preem.service
[Unit]
Description=Start preemtibles vm

[Service]
ExecStart=/usr/bin/start_preem
Restart=on-failure

[Install]
WantedBy=multi-user.target
