# real lab/start_prem.service
[Unit]
Description=Start preemtibles vm

[Service]
ExecStart=/usr/bin/start_prem
Restart=on-failure

[Install]
WantedBy=multi-user.target
